# omada-controller-runenv.target
[Unit]
Description=Runtime environment for the Omada EAP controller
PartOf=omada-controller.service
